# astrofin-sentinel-platform/infrastructure/asurdev/ml_engine/inference/ml-inference.service
[Unit]
Description=ML Inference API (XGBoost risk scoring)
After=network-online.target
Wants=network-online.target

[Service]
User=acos
Group=acos
WorkingDirectory=/home/acos/repos/AsurDev
Environment="ML_MODEL_PATH=/home/acos/repos/AsurDev/models/failure_xgb_v2.pkl"
Environment="ML_FEATURES_PATH=/home/acos/repos/AsurDev/models/features.txt"
Environment="PYTHONPATH=/home/acos/repos/AsurDev"

ExecStart=/home/acos/.venv/bin/uvicorn ml_engine.inference.api:app \
    --host 0.0.0.0 \
    --port 8081 \
    --access-log \
    --no-use-colors

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ml-inference

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/acos/repos/AsurDev/logs
PrivateTmp=true
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
